# display/kiosk.service
[Unit]
Description=Digital Signage Kiosk
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pedro
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pedro/.Xauthority
WorkingDirectory=/home/pedro/kiosk/display
ExecStart=/bin/bash /home/pedro/kiosk/display/start_kiosk.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target 